# deploy/hestia-api.container
[Unit]
Description=Hestia API Server
After=hestia-postgres.service

[Container]
ContainerName=hestia-api
Image=ghcr.io/zighome24/hestia-api:latest
Network=hestia.network
EnvironmentFile=/var/lib/hestia/hestia.env
Volume=/var/lib/hestia/uploads:/uploads:Z
Environment=PORT=8080
Environment=STORAGE_PATH=/uploads
HealthCmd=wget -q --spider http://localhost:8080/api/health || exit 1
HealthInterval=15s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target default.target
